# studio/deploy/systemd-user/kai-engine.service
[Unit]
Description=Kai support engine (uvicorn :6090)
PartOf=kai.target
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/workspace/kai
EnvironmentFile=-%h/.config/kai/kai-stack.env
Environment=PYTHONPATH=%h/workspace/kai
ExecStartPre=/usr/bin/test -f %h/workspace/kai/app.py
ExecStart=/bin/bash -c 'exec "${KAI_PYTHON:?Set KAI_PYTHON in kai-stack.env}" -m uvicorn app:app --host 127.0.0.1 --port 6090'
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=kai.target
